# post_install.sh
#!/bin/sh
. /etc/rc.subr
name="movienight"

# Start the service
/usr/local/etc/rc.d/${name} start